# files/etc/uci-defaults/98-set-root-password
#!/bin/sh

ROOT_HASH='$6$.piwIAoNrsjPMutX$1v6Ei.3N5S1P/Td4uSshl9w1Mya/zZT1.MzQ/8zvEIDegkXuo94LrN9md4JlFYRNHWcZUznFVH3ZGEC3XToYa0'

sed -i "s|^root:[^:]*:|root:${ROOT_HASH}:|" /etc/shadow

exit 0
