# files/etc/uci-defaults/99-set-defaults
#!/bin/sh

uci set system.@system[0].timezone='CST-8'
uci set system.@system[0].zonename='Asia/Shanghai'
uci commit system

uci set luci.diag.dns='www.qq.com'
uci set luci.diag.ping='www.qq.com'
uci set luci.diag.route='www.qq.com'
uci commit luci

uci set dhcp.@dnsmasq[0].cachesize='0'
uci set dhcp.@dnsmasq[0].min_cache_ttl='0'
uci set dhcp.@dnsmasq[0].noresolv='1'

uci set dhcp.lan.ignore='1'

uci -q del dhcp.wan
uci -q del dhcp.lan.ra
uci -q del dhcp.lan.ra_slaac
uci -q del dhcp.lan.ra_flags
uci -q del dhcp.lan.dhcpv6
uci -q del dhcp.lan.ndp
uci commit dhcp
